# Pi/PyBlaster2/cmake/debian/prerm
#!/bin/sh


rm /opt/PyBlaster2/python/*.pyc

update-rc.d pyblaster2 remove
