# Pi/PyBlaster2/cmake/debian/prerm
#!/bin/sh


rm /opt/PyBlaster2/python/*.pyc

rm /usr/bin/pyblaster2

update-rc.d pyblaster2 remove
